# LUB/crontab
#! /bin/bash

00 18 * * 1-5 /home/blox_land/LUB/LUB.sh

00 18 * * 6 /home/blox_land/LUB/SATURDAY_LUB.sh

00 18 * * 0 /home/blox_land/LUB/SUNDAY_LUB.sh

# $ crontab -e
